# r5s/files/etc/crontabs/root
0 4 * * * $HOME/cron/daily.sh &>/dev/null
0 6 1 * * $HOME/cron/monthly.sh &>/dev/null
